# rpi/blindstick.service
[Unit]
Description=Smart Blind Stick Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/blindstick
ExecStart=/usr/bin/python3 /home/pi/blindstick/main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
